# webcams/crontab
* * * * * tsocks /tmp/info-beamer/samples/webcams/cron.sh
